# overlay/etc/services.d/dropbox/finish
#!/bin/bash

rm -f /dbox/.dropbox/dropbox.pid
exit 0
