# overlay/etc/services.d/dropbox/run
#!/usr/bin/with-contenv bash

echo " - starting Dropbox daemon..."
s6-envuidgid -u dbox /dbox/.dropbox-dist/dropboxd
echo " - done."
