# overlay/etc/services.d/dropbox/run
#!/bin/bash

echo " - running Dropbox daemon..."
rm -f /dbox/.dropbox/dropbox.pid
#su - dbox -c "/dbox/.dropbox-dist/dropboxd"
runuser -l -u dbox /dbox/.dropbox-dist/dropboxd
