# overlay/etc/services.d/dropbox/run
#!/bin/bash

/dbox/.dropbox-dist/dropboxd
