# overlay/etc/services.d/dropbox/run
#!/bin/bash

echo " - starting Dropbox daemon..."
su -u dbox
~/.dropbox-dist/dropboxd
echo " - done."
