# overlay/etc/services.d/dropbox/run
#!/usr/bin/with-contenv bash

echo " - starting Dropbox daemon..."
s6-setuidgid $UID:$GID "/dbox/cli start"
echo " - done."
