# overlay/etc/services.d/dropbox/run
#!/bin/bash

echo " - starting Dropbox daemon..."
su - dbox -c "~/.dropbox-dist/dropboxd"
