# overlay/etc/services.d/dropbox/run
#!/bin/bash

echo " - running Dropbox daemon..."
rm -f /dbox/.dropbox/dropbox.pid

su - dbox -c "/etc/services.d/dropbox/status &"
su - dbox -c "/dbox/.dropbox-dist/dropboxd"
