# overlay/etc/services.d/dropbox/run
#!/usr/bin/with-contenv bash

echo " - starting Dropbox daemon..."
s6-setuidgid -B $UID:$GID /dbox/.dropbox-dist/dropboxd
echo " - done."
