# overlay/etc/services.d/dropbox/run
#!/bin/bash

STATUS=$(su - dbox -c "~/cli start && ~/cli status")
echo "$STATUS"
