# overlay/etc/services.d/dropbox/run
#!/bin/bash

echo " - starting Dropbox daemon..."
STATUS=$(su - dbox -c "~/cli running")
if [ "$STATUS" == "0" ] ; then
    su - dbox -c "~/.dropbox-dist/dropboxd"
fi
