# overlay/etc/services.d/dropbox/run
#!/bin/bash

echo " - Dropbox daemon..."
rm -f /dbox/.dropbox/dropbox.pid
runuser dbox -c "~/.dropbox-dist/dropboxd"
