# overlay/etc/services.d/dropbox/run
#!/bin/bash

echo " - starting Dropbox daemon..."
rm -f /dbox/.dropbox/dropbox.pid
su - dbox -c "~/.dropbox-dist/dropboxd"
