# overlay/etc/services.d/dropbox/run
#!/bin/bash

su - dbox -c "~/cli start -i"
